# client/dist-specific/deb/prerm
#! /usr/bin/sh
systemctl stop turn-proxy-server || true
systemctl disable turn-proxy-server || true